# systemd/vps1777-secrets-check.timer
# vps1777-secrets-check.timer — check settimanale delle scadenze secret.
#
# I secret invecchiano lentamente → cadenza settimanale (una nudge/settimana se
# qualcosa è scaduto). RandomizedDelaySec distribuisce il carico; Persistent
# recupera i check persi a VPS spenta.

[Unit]
Description=vps1777 — check settimanale scadenze secret

[Timer]
OnCalendar=weekly
RandomizedDelaySec=6h
Persistent=true

[Install]
WantedBy=timers.target
