# systemd/vps1777-secrets-check.service
# vps1777-secrets-check.service — età dei secret + notifica scadenze.
#
# Scrive onboarding/secrets_status.json (letto dalla pagina /admin/secrets) e
# notifica il owner su Telegram i secret oltre soglia. Legge l'mtime dei file in
# secrets/ e — best-effort — la freschezza dei cookie NotebookLM nel volume
# nlm-auth (H37, via `docker run` read-only). Zero rete se non per Telegram.
#
# H43 (il bug che spegneva il monitoraggio): utente e path NON sono hardcodati.
# La CLI (install_systemd_units → render_unit) sostituisce @OPERATOR_USER@ /
# @REPO@ con i valori REALI a ogni update; il deploy iniziale fa lo stesso (vedi
# note per deploy.sh). Con un operatore diverso da "vps1777" le scadenze dei
# secret NON smettono più di essere controllate in silenzio.

[Unit]
Description=vps1777 — check scadenze secret
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
User=@OPERATOR_USER@
Group=@OPERATOR_USER@
Environment=VPS1777_HOME=@REPO@
WorkingDirectory=@REPO@
ExecStart=/usr/local/bin/vps1777 secrets-status --notify
TimeoutStartSec=120
# H43 hardening: come check-update. Niente sudo, nessuna scrittura fuori dal
# repo. Il probe dei cookie nlm usa `docker run` (socket in /run, non /home) →
# ProtectHome=read-only non lo tocca; la scrittura resta confinata al repo.
NoNewPrivileges=true
ProtectHome=read-only
ReadWritePaths=@REPO@
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
LockPersonality=true
